# dist/com.github.paveq.ebeco-spot.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template for the ebeco-spot LaunchAgent. `make install` substitutes the
  __PLACEHOLDERS__ with absolute paths and writes the result to
  ~/Library/LaunchAgents/. Do not edit the installed copy by hand; re-run
  `make install` instead.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>__LABEL__</string>

	<!--
	  Run the binary directly. -keychain makes it read the credentials from the
	  Keychain itself (so the item's ACL can be bound to this binary, not to
	  /usr/bin/security); -log oslog routes logs to the unified logging system.
	-->
	<key>ProgramArguments</key>
	<array>
		<string>__PREFIX__/ebeco-spot</string>
		<string>-config</string>
		<string>__PREFIX__/config.toml</string>
		<string>-keychain</string>
		<string>-log</string>
		<string>oslog</string>
	</array>

	<!-- Start at login and keep it alive; launchd restarts it if it exits. -->
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>30</integer>

	<!-- Low-priority background work (does not keep the system awake). -->
	<key>ProcessType</key>
	<string>Background</string>

	<!--
	  No StandardOutPath/StandardErrorPath on purpose: run.sh passes -log oslog,
	  so the binary logs into the unified logging system (os_log), which the OS
	  caps and rotates automatically. View it with:
	    log stream --predicate 'subsystem == "com.github.paveq.ebeco-spot"'   (live)
	    log show   --predicate 'subsystem == "com.github.paveq.ebeco-spot"' --last 1h
	-->
</dict>
</plist>
