# dist/com.github.paveq.ebeco-spot.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template for the ebeco-spot LaunchAgent. `make install` substitutes the
  __PLACEHOLDERS__ with absolute paths and writes the result to
  ~/Library/LaunchAgents/. Do not edit the installed copy by hand; re-run
  `make install` instead.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>__LABEL__</string>

	<!-- run.sh reads credentials from the Keychain, then execs the binary. -->
	<key>ProgramArguments</key>
	<array>
		<string>__PREFIX__/run.sh</string>
	</array>

	<!-- Start at login and keep it alive; launchd restarts it if it exits. -->
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>30</integer>

	<!-- Low-priority background work (does not keep the system awake). -->
	<key>ProcessType</key>
	<string>Background</string>

	<key>StandardOutPath</key>
	<string>__LOGFILE__</string>
	<key>StandardErrorPath</key>
	<string>__LOGFILE__</string>
</dict>
</plist>
